# birdhouse/components/proxy/cron.d/logrotate
* * * * * root /usr/sbin/logrotate -v /etc/logrotate.d/nginx
